# package/rpm/postinst
#!/bin/bash

set +e
chkconfig --add pfsd_env
set -e
